# base-container/conf/etc/supervisor.d/watchkeys.conf
[group:watchkeys]
programs=watchkeysd
priority=999

[program:watchkeysd]
command = /opt/docker/bin/service.d/watchkeys.sh
autostart = false
autorestart = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0